# cdk/assets-ec2/postgres-backup-s3/src/crontab
* * * * * /bin/sh /src/backup.sh >> /var/log/cron.log 2>&1